# roles/daily-backup/files/daily-backup.service
[Unit]
Description=Full server backup

[Service]
Type=simple
ExecStart=/root/fullBackup.sh

